# systemd/user/backups.service
[Unit]
Description=Backup Service

[Service]
ExecStart=/usr/local/bin/backups
