# systemd/user/compton.service
[Unit]
Description=Compton Service

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/compton

[Install]
WantedBy=default.target
